# ops/mystmon.cron
0 */6 * * * cd /mnt/ssd/codex/mystmon && /usr/bin/docker compose run --rm mystmon python -c "import asyncio; from mystmon.api import create_app; app=create_app(); asyncio.run(app.state.scheduler.collect_once())"
